# files/etc/systemd/system/update-firefox.timer

[Unit]
Description=Automatically update the local firefox installation

[Timer]
OnCalendar=06:00

[Install]
WantedBy=timers.target

