# files/desktop/etc/systemd/system/network.mounts.automount.j2

[Unit]
Description=Automount {{ path }}

[Automount]
Where={{ path }}
{% if idle_timeout != 0 %}
TimeoutIdleSec={{ idle_timeout }}
{% endif %}

[Install]
WantedBy=multi-user.target

