# files/home/kmchan/.config/systemd/user/bumblebee-proxy.service

[Unit]
After=bumblebeed.service

[Service]
ExecStart=/lib/systemd/systemd-socket-proxyd /var/run/bumblebee.socket
PrivateTmp=yes

