# files/home/kmchan/.config/systemd/user/udisks2temp.service

[Unit]
Description=HDDTemp Bridge for UDisks2

[Service]
ExecStart=/home/kmchan/bin/udisks2temp

[Install]
WantedBy=default.target

