# files/etc/profile.d/video.sh.j2

PATH=${PATH}:/opt/prime-run/bin
export PATH

