# files/desktop/etc/profile.d/windows.sh.j2

if id -nGz | grep -qzxF windows ; then
  PATH=${PATH}:/opt/windows/bin
fi

